# data/fleet-commander-clientad.service.in
[Unit]
Description=Fleet Commander Client AD dbus service

[Service]
Type=dbus
BusName=org.freedesktop.FleetCommanderClientAD
Environment=PYTHONPATH=@FCPYTHONDIR@
ExecStart=@PYTHON@ -m fleetcommanderclient.fcclientad --configuration @XDGCONFIGDIR@/fleet-commander-client.conf
StandardOutput=syslog
StandardError=inherit

[Install]
WantedBy=multi-user.target
