# data/fleet-commander-adclient.service.in
[Unit]
Description=Fleet Commander Client Active Directory service

[Service]
Type=simple
Environment=PYTHONPATH=@FCPYTHONDIR@
ExecStart=@PYTHON@ -m fleetcommanderclient.fcadclient --configuration @XDGCONFIGDIR@/fleet-commander-client.conf
StandardOutput=syslog
StandardError=inherit

[Install]
WantedBy=default.target
