# data/fleet-commander-client.service.in
[Unit]
Description=Fleet Commander Client dbus service

[Service]
Type=dbus
BusName=org.freedesktop.FleetCommanderClient
Environment=PYTHONPATH=@FCPYTHONDIR@
ExecStart=@PYTHON@ -m fleetcommander.fcclient --configuration @XDGCONFIGDIR@/fleet-commander-client.conf
StandardOutput=syslog
StandardError=inherit

[Install]
WantedBy=multi-user.target
